# deploy/systemd/snapshot-gc.service
[Unit]
Description=TianGong LCA Snapshot Storage GC
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/projects/lca_workspace/tiangong-lca-calculator
EnvironmentFile=/home/ubuntu/projects/lca_workspace/tiangong-lca-calculator/.env
Environment=RUST_LOG=info
Environment=MAINTENANCE_JOB_ENVIRONMENT=main
SyslogIdentifier=maintenance_enqueue
ExecStart=/home/ubuntu/projects/lca_workspace/tiangong-lca-calculator/target/release/maintenance_enqueue snapshot-gc --execute --snapshot-retention-days 30 --orphan-retention-days 30 --max-snapshots 100 --max-orphan-dirs 200 --max-bytes 2147483648 --batch-size 50
